# systemd/user/podman-image-prune.timer
[Unit]
Description=Delete unused podman images periodically

[Timer]
OnCalendar=Mon *-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
