# systemd/user/podman-image-prune.service
[Unit]
Description=Delete unused podman images

[Service]
Type=oneshot
ExecStart=podman image prune -af
