# systemd/user/quadlet-sync.service
[Unit]
Description=Sync quadlets from git

[Service]
Type=oneshot
ExecStart=quadlet-sync
Environment=GIT_URL=
Environment=GIT_SUBDIR=
