# quadlet-sync/quadlet-sync.service
[Unit]
Description=Sync quadlets from git

[Service]
Type=oneshot
ExecStart=quadlet-sync
